# luci-theme-fibernet/root/etc/uci-defaults/30_luci-theme-fibernet
#!/bin/sh
uci batch <<-EOF
	set luci.themes.Fibernet=/luci-static/fibernet
	set luci.main.mediaurlbase=/luci-static/fibernet
	commit luci
EOF
exit 0
